# webby/rc.local
#!/bin/sh

touch /home/arducam/wifiUP.txt
sleep 1
rm /home/arducam/wifiUP.txt
echo "Cleanup of wifi tracking file."

nohup /home/arducam/bin/webInterface.py >> /home/arducam/weblog.txt 2>&1 

exit 0
